# ansible/roles/configure_ttrss/templates/ttrss_feed_updater.service.j2
[Unit]
Description=ttrss_feed_updater
After=network.target mysql.service

[Service]
User=nginx
ExecStart={{ rss_hostname }}/update_daemon2.php

[Install]
WantedBy=multi-user.target